# app/service/s2s-api.conf
description "S2S Server"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python3 /vagrant/api/Stock2ShopAPI.py